# config/launchd/com.user.email-feedback.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.email-feedback</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jake/Developer/email-assistant/venv/bin/python</string>
        <string>/Users/jake/Developer/email-assistant/scripts/collect_feedback.py</string>
        <string>--days</string>
        <string>1</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jake/Developer/email-assistant/logs/feedback.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jake/Developer/email-assistant/logs/feedback-error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/jake/Developer/email-assistant</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
